# deploy/com.bosphorify.studiod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  studiod launchd system daemon plist (template).

  This file is a TEMPLATE: the literal string __TAILSCALE_IP__ below is
  substituted with the host's Tailscale IPv4 by deploy/install-server.sh
  before the rendered file is dropped into /Library/LaunchDaemons/.

  The daemon runs as root (UserName=root) so it can call powermetrics for
  GPU and power telemetry, which is otherwise a sudo-gated operation.
  Keep the dependency surface minimal: see pyproject.toml [collector].
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bosphorify.studiod</string>

    <key>UserName</key>
    <string>root</string>

    <key>GroupName</key>
    <string>wheel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/studiod/venv/bin/python</string>
        <string>-m</string>
        <string>collector</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>STUDIOD_BIND_HOST</key>
        <string>__TAILSCALE_IP__</string>
        <key>STUDIOD_BIND_PORT</key>
        <string>8765</string>
        <key>STUDIOD_TOKEN_FILE</key>
        <string>/etc/studiod/token</string>
        <key>STUDIOD_TMUX_USER</key>
        <string>__TMUX_USER__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>C</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>/var/log/studiod.out.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/studiod.err.log</string>

    <!-- Throttle restart attempts so a crash loop does not hammer the
         system. launchd defaults to 10s; we bump it to 30 to give us a
         chance to ssh in and look at logs. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
